# applications/welcome/cron/crontab
#crontab
*5 24  * *  * root *applications/instore/modules/updater.py
